# src/etc/profile.d/pmc-env.sh
XDEBUG_REMOTE_HOST=10.0.2.2
PMC_PHPUNIT_BOOTSTRAP=/srv/www/pmc/wp-config/bootstrap-phpunit.php
WP_TESTS_DIR=/srv/www/wp-core/wp-tests/tests/phpunit/
